# config/systemd/user/docker-prune.timer
[Unit]
Description=Run docker prune in intervals

[Timer]
OnStartupSec=1min
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

